# root/etc/services.d/mongodb/run
#!/usr/bin/with-contenv bash

exec \
        s6-setuidgid abc /usr/bin/mongod -f /config/mongodb.conf --logpath /config/log/mongodb.log